# Raspberry Pi/CPU Temperature monitor/crontab
@reboot python3 /home/username/cpu_temp.py
